# recipes-devtools/supervisor/files/configs.d/cantransceiver.conf
[program:CanTransceiver]
command=/usr/bin/CanTransceiver
autostart=true
autorestart=true
# stderr_logfile=/var/log/CanReceiver.err.log
# stdout_logfile=/var/log/CanReceiver.out.log
environment=LD_LIBRARY_PATH="/usr/lib"

