# recipes-core/network-setup/files/network-setup.service
[Unit]
Description=Custom Network Setup
After=network.target

[Service]
ExecStart=/usr/bin/network-setup.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
